# systemd/kudos-backup.timer
[Unit]
Description=Daily backup of Kudos database

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target